# deploy/systemd/nifty-eod-retry.service
[Unit]
Description=NIFTY NSE EOD retry-missing
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=__USER__
WorkingDirectory=__APPDIR__
Environment=TZ=Asia/Kolkata
ExecStart=__APPDIR__/.venv/bin/python -m nifty.jobs eod --retry-missing
