# deploy/systemd/nifty-dashboard-stop.timer
[Unit]
Description=Stop dashboard 16:05 IST

[Timer]
OnCalendar=Mon..Fri 16:05 Asia/Kolkata
Persistent=false

[Install]
WantedBy=timers.target
