# deploy/systemd/nifty-gift-start.timer
[Unit]
Description=Start GIFT monitor 15:35 IST

[Timer]
OnCalendar=Mon..Fri 15:35 Asia/Kolkata
Persistent=false

[Install]
WantedBy=timers.target
