# deploy/systemd/nifty-eod-cm.timer
[Unit]
Description=NSE EOD cash reports 18:30 IST

[Timer]
OnCalendar=Mon..Fri 18:30 Asia/Kolkata
Persistent=false

[Install]
WantedBy=timers.target
