# deploy/systemd/nifty-dashboard-start.service
[Unit]
Description=Start nifty-dashboard.service

[Service]
Type=oneshot
ExecStart=/bin/systemctl start nifty-dashboard.service
