# deploy/systemd/nifty-eod-fo.service
[Unit]
Description=NIFTY NSE EOD download (fo)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=__USER__
WorkingDirectory=__APPDIR__
Environment=TZ=Asia/Kolkata
ExecStart=__APPDIR__/.venv/bin/python -m nifty.jobs eod --targets fo_bhavcopy,participant_oi,participant_vol
