# deploy/systemd/nifty-gift-start.service
[Unit]
Description=Start nifty-gift.service

[Service]
Type=oneshot
ExecStart=/bin/systemctl start nifty-gift.service
