# deploy/systemd/nifty-dashboard-stop.service
[Unit]
Description=Stop nifty-dashboard.service

[Service]
Type=oneshot
ExecStart=/bin/systemctl stop nifty-dashboard.service
