# deploy/systemd/nifty-gift-stop.service
[Unit]
Description=Stop nifty-gift.service

[Service]
Type=oneshot
ExecStart=/bin/systemctl stop nifty-gift.service
