# deploy/systemd/nifty-dashboard.service
[Unit]
Description=NIFTY OI Velocity Dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__APPDIR__
Environment=TZ=Asia/Kolkata
ExecStart=__APPDIR__/.venv/bin/python -m nifty.dashboard --host 127.0.0.1 --port 8080
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
