# deploy/systemd/nifty-dashboard-start.timer
[Unit]
Description=Start dashboard 08:00 IST

[Timer]
OnCalendar=Mon..Fri 08:00 Asia/Kolkata
Persistent=false

[Install]
WantedBy=timers.target
