# deploy/systemd/nifty-morning.timer
[Unit]
Description=Morning pipeline 08:15 IST

[Timer]
OnCalendar=Mon..Fri 08:15 Asia/Kolkata
Persistent=false
Unit=nifty-job@morning.service

[Install]
WantedBy=timers.target
